# apps/backend/scripts/release/cohesion-package.service
[Unit]
Description=Cohesion file service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=cohesion
Group=cohesion
WorkingDirectory=/var/lib/cohesion
Environment=HOME=/var/lib/cohesion
Environment=COHESION_RUNTIME_ROOT=/var/lib/cohesion/runtime
Environment=COHESION_INSTALL_CHANNEL=package
Environment=COHESION_LAUNCH_MODE=background
ExecStart=/usr/bin/cohesion
Restart=on-failure
RestartSec=2s
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
